# deploy/web-chat.service.template
# Шаблон systemd unit для web-chat.
# Генерируется скриптом deploy/install.sh — не редактируйте пути вручную в /etc/systemd/system/.
# Плейсхолдеры: @@INSTALL_ROOT@@, @@RUN_USER@@, @@RUN_GROUP@@, @@WEB_PORT@@

[Unit]
Description=web-chat — LAN-чат с AI-агентом, MCP и Stable Diffusion
Documentation=file://@@INSTALL_ROOT@@/deploy/DEPLOY.md
After=network-online.target postgresql.service
Wants=network-online.target
# postgresql.service — при DATABASE_URL=postgresql; если Postgres не установлен, unit игнорируется

[Service]
Type=simple
User=@@RUN_USER@@
Group=@@RUN_GROUP@@
WorkingDirectory=@@INSTALL_ROOT@@
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=@@INSTALL_ROOT@@/.env
ExecStart=@@INSTALL_ROOT@@/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port @@WEB_PORT@@
Restart=on-failure
RestartSec=5
TimeoutStopSec=60

# Долгие запросы к SD / LLM
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
